# .cron
* * * * * PATH="/usr/local/bin:${PATH}" ${HOME}/.crontab/ps1-all
